# callchecker-analysis.service
[Unit]
Description=Callchecker Dialog Analysis Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=callchecker
Group=callchecker
WorkingDirectory=/home/callchecker/callchecker
Environment=PATH=/home/callchecker/callchecker/venv/bin
ExecStart=/home/callchecker/callchecker/venv/bin/python dialog_analysis.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
