# callchecker-bitrix.service
[Unit]
Description=Callchecker Bitrix24 Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=callchecker
Group=callchecker
WorkingDirectory=/home/callchecker/callchecker
Environment=PATH=/home/callchecker/callchecker/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/home/callchecker/callchecker/venv/bin/python bitrix24/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
